# ha_media_request/rootfs/etc/s6-overlay/s6-rc.d/media-requests/run
#!/usr/bin/with-contenv bashio

exec python3 /app/media_requests.py
