# plex_wol_listener/rootfs/etc/s6-overlay/s6-rc.d/plex-wol/finish
#!/usr/bin/with-contenv bashio
bashio::log.warning "Plex WoL proxy exited (code: $1, signal: $2). Restarting in 5 seconds…"
sleep 5
