# plex_wol_listener/rootfs/etc/s6-overlay/s6-rc.d/plex-wol/run
#!/usr/bin/with-contenv bashio
python3 -u /migrate.py
exec python3 -u /listener.py
